# systemd/overrides/vllm.service.d/10-stack.conf
[Unit]
PartOf=rag-stack.target
After=network-online.target
Wants=network-online.target

