# systemd/rag-gateway.service
[Unit]
Description=RAG Gateway (FastAPI)
After=network-online.target qdrant.service tei-embed.service tei-rerank.service vllm.service
Wants=network-online.target
Requires=qdrant.service tei-embed.service tei-rerank.service vllm.service

[Service]
Type=simple
User=rag
Group=rag
WorkingDirectory=/opt/llm/rag-gateway

# Wait for Qdrant
ExecStartPre=/opt/llm/rag-gateway/bin/wait-http.sh --name qdrant --url http://127.0.0.1:6333/collections --method GET --timeout 60

# Wait for TEI embed (POST) - use bash so JSON is always preserved correctly
ExecStartPre=/usr/bin/bash -lc '/opt/llm/rag-gateway/bin/wait-http.sh --name tei-embed --url http://127.0.0.1:8081/v1/embeddings --method POST --json "{\"model\":\"Qwen/Qwen3-Embedding-8B\",\"input\":\"ready\"}" --timeout 120'

# Wait for vLLM to publish models; we just need to see "generator" in /v1/models
ExecStartPre=/usr/bin/bash -lc '/opt/llm/rag-gateway/bin/wait-http.sh --name vllm --url http://127.0.0.1:8000/v1/models --method GET --expect "\"id\"[[:space:]]*:[[:space:]]*\"generator\"" --timeout 180'

Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/llm/rag-gateway/.venv/bin/python -m rag_gateway.main --host 127.0.0.1 --port 9000

Restart=on-failure
RestartSec=2

[Install]
WantedBy=rag-stack.target

