# systemd/rag-gateway.service
[Unit]
Description=RAG Gateway (FastAPI)
After=network-online.target qdrant.service tei-embed.service tei-rerank.service vllm.service
Wants=network-online.target
Requires=qdrant.service tei-embed.service tei-rerank.service vllm.service

[Service]
Type=simple
User=rag
Group=rag
WorkingDirectory=/opt/llm/rag-gateway

Environment=RAG_GATEWAY_CONFIG=/etc/rag-gateway/config.yaml

# Readiness gating (order != readiness). Use helper script installed by install.sh.
ExecStartPre=/opt/llm/rag-gateway/bin/wait-http.sh --name qdrant --url http://127.0.0.1:6333/collections --method GET --timeout 180

# vLLM: don't start gateway until vLLM is answering OpenAI-compatible requests
ExecStartPre=/opt/llm/rag-gateway/bin/wait-http.sh --name vllm --url http://127.0.0.1:8000/v1/models --method GET --timeout 600

# TEI embed: send valid JSON as ONE argument; use list input to avoid 4xx
ExecStartPre=/opt/llm/rag-gateway/bin/wait-http.sh --name tei-embed --url http://127.0.0.1:8081/v1/embeddings --method POST --json '{"model":"Qwen/Qwen3-Embedding-8B","input":["ready"]}' --timeout 600

# TEI rerank: match the schema your client uses (texts)
ExecStartPre=/opt/llm/rag-gateway/bin/wait-http.sh --name tei-rerank --url http://127.0.0.1:8082/v1/rerank --method POST --json '{"model":"BAAI/bge-reranker-large","query":"ready","texts":["ready"]}' --timeout 600


ExecStart=/opt/llm/rag-gateway/.venv/bin/uvicorn rag_gateway.app:app --host 127.0.0.1 --port 9000

Restart=always
RestartSec=5
TimeoutStartSec=20min

[Install]
WantedBy=multi-user.target

