# systemd/rag-gateway.service
[Unit]
Description=RAG Gateway (FastAPI)
After=network.target qdrant.service tei-embed.service tei-rerank.service vllm.service
Wants=qdrant.service tei-embed.service tei-rerank.service vllm.service

[Service]
Type=simple
User=rag
Group=rag
WorkingDirectory=/opt/llm/rag-gateway

# Wait for dependencies to be reachable before starting the API
ExecStartPre=/opt/llm/rag-gateway/bin/wait-http.sh --name qdrant --url http://127.0.0.1:6333/collections --method GET --timeout 60

# IMPORTANT: escape JSON quotes so systemd does not strip them
ExecStartPre=/opt/llm/rag-gateway/bin/wait-http.sh --name tei-embed --url http://127.0.0.1:8081/v1/embeddings --method POST --timeout 180 --json "{\"model\":\"Qwen/Qwen3-Embedding-8B\",\"input\":\"ready\"}"

ExecStartPre=/opt/llm/rag-gateway/bin/wait-http.sh --name vllm --url http://127.0.0.1:8000/v1/models --method GET --timeout 180

ExecStart=/opt/llm/rag-gateway/.venv/bin/uvicorn rag_gateway.app:app --host 0.0.0.0 --port 8001

Restart=always
RestartSec=2
TimeoutStartSec=0

# Reasonable defaults
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target

