# systemd/rag-crawl.timer
[Unit]
Description=RAG Scheduled Crawl Timer (ticks every minute; interval in config)

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=30s
Unit=rag-crawl.service

[Install]
WantedBy=timers.target

