# systemd/rag-crawl.service
[Unit]
Description=RAG Scheduled Crawl (incremental)
After=network.target qdrant.service tei-embed.service tei-rerank.service
Wants=qdrant.service tei-embed.service tei-rerank.service

[Service]
Type=oneshot
User=rag
Group=rag
WorkingDirectory=/opt/llm/rag-gateway
ExecStart=/opt/llm/rag-gateway/.venv/bin/python -m rag_gateway.scheduled_crawl_cli --config /etc/rag-gateway/config.yaml

