# systemd/overrides/tei-rerank.service.d/10-stack.conf
[Unit]
PartOf=rag-stack.target
After=network-online.target qdrant.service tei-embed.service
Wants=network-online.target
Requires=qdrant.service tei-embed.service

