# systemd/overrides/tei-embed.service.d/10-stack.conf
[Unit]
PartOf=rag-stack.target
After=network-online.target qdrant.service
Wants=network-online.target
Requires=qdrant.service

