# systemd/overrides/rag-crawl.timer.d/10-stack.conf
[Unit]
PartOf=rag-stack.target
After=rag-gateway.service
Requires=rag-gateway.service

